# mlflow.ini
[program:mlflow]
command=sh /data/mlflow/examples/run_mlflow.sh
user=root
priority=999
startsecs=10
redirect_stderr=True
autostart=True
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=/data/mlflow/examples/logout.out